# infra/deploy/vm/systemd/margadeep-mcp.service
[Unit]
Description=Margadeep MCP server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/margadeep/margadeep.env
ExecStart=/usr/local/bin/margadeep-run-mcp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
